# systemd/teebot-poller.timer
[Unit]
Description=Run TeeBot poller every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
